# systemd/hybrid-memory-maintenance.timer
[Unit]
Description=Run Memory Engine hybrid maintenance daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
